# etc/profile.d/profile.sh
echo "sourcing profile.sh"
alias ll="ls -alh"
alias claude-yolo="claude --dangerously-skip-permissions"
